# src/telemetrics/mock-failing-service.service
[Unit]
Description=Mock Failing Service

[Service]
ExecStart=/usr/bin/false

[Install]
WantedBy=multi-user.target
