# src/telemetrics/mock-telemetrics.service
[Unit]
Description=Mock Telemetrics Backend

[Service]
ExecStart=/etc/telemetrics/mock-telemetrics-backend.sh

[Install]
WantedBy=multi-user.target
